# gormes/scripts/orchestrator/systemd/gormes-orchestrator-audit.timer.in
[Unit]
Description=Run gormes-orchestrator-audit every 20 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=20min
AccuracySec=30s
Persistent=true
Unit=gormes-orchestrator-audit.service

[Install]
WantedBy=timers.target
